# rme-mqtt-bridge.service
[Unit]
Description=RME ADI-2 DAC MQTT <-> MIDI Bridge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/default/rme-mqtt-bridge
Environment=MQTT_HOST=192.168.100.60
Environment=MQTT_PORT=1883
Environment=MIDI_PORT=hw:1,0,0
Environment=READY_STREAK=5
Environment=APPLY_RETRIES=10
Environment=DEBUG=0
Environment=APPLY_RETRY_DELAY=0.8
ExecStart=/usr/local/bin/rme_mqtt_bridge.py
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
